# local-runner/com.mountain.inference.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd template for the local Mac inference runner.

  1. Replace __PROJECT_PATH__ with the absolute path to your checkout
     (e.g. /Users/tommydoerr/dev/is-the-mountain-out).
  2. Copy to ~/Library/LaunchAgents/com.mountain.inference.plist
  3. launchctl load ~/Library/LaunchAgents/com.mountain.inference.plist

  Tail the log:
    tail -f /tmp/mountain_inference.out /tmp/mountain_inference.err

  Unload:
    launchctl unload ~/Library/LaunchAgents/com.mountain.inference.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mountain.inference</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_PATH__/local-runner/tick.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>__PROJECT_PATH__</string>

    <key>StandardOutPath</key>
    <string>/tmp/mountain_inference.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mountain_inference.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
